# core/systemd_example.service
[Unit]
Description=Backhaul Watchdog Service
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash /root/backhaul/backhaul_watchdog.sh run
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
